# apache/root/etc/s6/apache/run
#!/usr/bin/env bash

if [[ -f ./setup ]]; then
  ./setup
fi

source /etc/apache2/envvars

exec apache2ctl -D FOREGROUND
