# elasticsearch/root/etc/s6/elasticsearch/finish
#!/usr/bin/env bash

s6-svscanctl -t /etc/s6
