# mongodb/root/etc/s6/mongodb/finish
#!/usr/bin/env bash

if [[ ! -s /var/lib/mongo/mongodb/mongod.lock ]]; then
  rm --verbose /var/lib/mongo/mongodb/mongod.lock
fi

s6-svscanctl -t /etc/s6
