# mongodb/root/etc/s6/mongodb/run
#!/usr/bin/env bash

mkdir -p /var/lib/mongo/mongodb
chown -R mongodb /var/lib/mongo/mongodb

exec s6-setuidgid mongodb mongod --config /etc/mongod.conf
