# logstash/root/etc/s6/rinetd/run
#!/usr/bin/env bash

exec rinetd -c /etc/rinetd.conf -f
