# base/root/etc/s6/logstash-forwarder/run
#!/usr/bin/env bash

if [[ -f ./setup ]]; then
  ./setup
fi

exec /opt/logstash-forwarder/bin/logstash-forwarder -config /etc/logstash-forwarder.conf > /dev/null 2>/dev/null
